# deploy/systemd/tsm-detector.service
[Unit]
Description=TSM Python AI Threat Detector
Documentation=https://github.com/tsm79/tsm79
After=network-online.target
Wants=network-online.target
PartOf=tsm.target

[Service]
Type=simple
User=tsm
Group=tsm

ExecStart=/usr/local/bin/tsm-detector
WorkingDirectory=/var/lib/tsm/detector

EnvironmentFile=/etc/tsm/detector.env
Environment=TSM_DETECTOR_PORT=8001
Environment=TSM_DETECTOR_WORKERS=4
Environment=PYTHONUNBUFFERED=1

Restart=on-failure
RestartSec=5s
StartLimitIntervalSec=120s
StartLimitBurst=3

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes

LimitNOFILE=16384
MemoryMax=1G

TimeoutStopSec=15s
KillMode=mixed

[Install]
WantedBy=multi-user.target tsm.target
